# silchenko/project_03/postgrest/postgrest.service
[Unit]
Description=PostgREST

[Service]
ExecStart=/bin/bash /opt/postgrest/run_postgrest.sh

[Install]
WantedBy=multi-user.target
